# airootfs/etc/systemd/system/run-tauri-sample.service
[Unit]
Description=Run a simple tauri app

[Service]
ExecStart=/usr/local/bin/run-sample-tauri-app


[Install]
WantedBy=multi-user.target
