# airootfs/etc/systemd/system/startx-onboot.service

[Unit]
Description=Script which starts X serve r

[Service]
ExecStart=startx


[Install]
WantedBy=multi-user.target
