# live-iso/airootfs/etc/systemd/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin difuse --noclear %I 38400 linux
